# netplanner/sriov/files/sriov-setup.service
[Unit]
Description=Configure SR-IOV VF functions on boot
DefaultDependencies=no
After=systemd-udevd.service
Before=network-pre.target openvswitch-switch.service mlnx-switchdev-mode.service
Wants=network-pre.target

[Service]
Type=oneshot
KillMode=none
ExecStart=/usr/local/sbin/netplanner --only-sriov configure
TimeoutSec=0

[Install]
WantedBy=multi-user.target
