# netplanner/sriov/files/sriov-udev.rules
ACTION=="add", SUBSYSTEM=="net", ATTRS{sriov_totalvfs}=="?*", RUN+="/usr/local/sbin/netplanner --only-sriov configure"